# .vagrant/vagrant-mounted.conf
start on vagrant-mounted

exec sudo service apache2 start
